# pi-player.service
[Unit]
Description=Pi Player
# network must be ready
After=network.target

[Service]
Type=simple
# if using custom user and group, define here.
#User=user
#Group=user
LimitNOFILE=1024

Restart=on-failure
RestartSec=10
StartLimitInterval=60

# change to your install directory.
WorkingDirectory=/home/user/dev/go/pi-player
ExecStart=/home/user/dev/go/pi-player/pi-player

[Install]
WantedBy=multi-user.target